# ccp/7_deploy_osh/roles/suse-osh-deploy/templates/helm-serve.service.j2
[Unit]
Description=Helm Server
After=network.target

[Service]
User=root
Restart=always
ExecStart=/usr/local/bin/helm serve

[Install]
WantedBy=multi-user.target
